# ai-router.service
[Unit]
Description=LoveSupremacy AI Router Pool
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/LoveSupremacy-Telegram-Bot
ExecStart=/usr/bin/node /opt/LoveSupremacy-Telegram-Bot/ai_router.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

Environment=NODE_ENV=production
Environment=OPENROUTER_API_KEY=${OPENROUTER_API_KEY}
Environment=SENSENOVA_API_KEY=${SENSENOVA_API_KEY}
Environment=GEMINI_API_KEY=${GEMINI_API_KEY}
Environment=SILICONFLOW_API_KEY=${SILICONFLOW_API_KEY}
Environment=EAS_SD_TOKEN=${EAS_SD_TOKEN}
Environment=EAS_SD_LORA_TRIGGER=${EAS_SD_LORA_TRIGGER}

[Install]
WantedBy=multi-user.target
