# system/nxsiran-bridge.service
[Unit]
Description=NxSiran Bridge Service (已弃用 - 功能已合并到 webhook_server.py)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/NxSiran-Telegram-Bot
EnvironmentFile=/opt/NxSiran-Telegram-Bot/.env
ExecStart=/usr/bin/python3 /opt/NxSiran-Telegram-Bot/system/webhook_server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
